# service/nova-common/run
#!/bin/sh

svcname=${PWD##*/}
exec runuser -u nova -g nova $svcname

